# moonserve.service
[Unit]
Description=Generate moon phase image
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=moonserve
ExecStart=/opt/moonserve/venv/bin/python /opt/moonserve/generate_moon.py
Environment=MOON_DATA=/opt/moonserve
Environment=MOON_OUTPUT=/var/www/moonserve/moon.png
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/moonserve /opt/moonserve
